# systemd/session-recall.service.example
[Unit]
Description=session-recall agent history backup (restic)
Documentation=https://github.com/restic/restic

[Service]
Type=oneshot
# Change these paths to your checkout.
ExecStart=%h/path/to/session-recall/backup.sh
StandardOutput=append:%h/path/to/session-recall/backup.log
StandardError=append:%h/path/to/session-recall/backup.log
# 跑慢一点也不抢资源（HDD 后端友好）
Nice=10
IOSchedulingClass=idle
